# .config/supervisor/uwsgi.conf
[program:uwsgi]
command = /root/.pyenv/versions/deploy_eb_docker/bin/uwsgi --ini /srv/deploy_eb_docker/.config/uwsgi/uwsgi-app.ini